# srcs/phpmyadmin/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = sh -c "nignx -g 'daemon off;' && kill -s SIGTERM $(cat supervisord.pid)"

[program:php-fpm]
command = sh -c "/usr/sbin/php-fpm7 -R -F && kill -s SIGTERM $(cat supervisord.pid)"
